# kamehouse-shell/bin/lin/pi/startup/pi-rc-local.service
[Unit]
After=network.service

[Service]
ExecStart=/home/KAMEHOUSE_USER/programs/kamehouse-shell/bin/lin/pi/startup/pi-rc-local.sh

[Install]
WantedBy=default.target
