# kamehouse-shell/bin/pi/startup/rc-local.service
[Unit]
After=network.service

[Service]
ExecStart=/home/KAMEHOUSE_USER/programs/kamehouse-shell/bin/pi/startup/rc-local.sh

[Install]
WantedBy=default.target
